# deploy/systemd/jarwo-gateway.service
[Unit]
Description=JARWO gateway + frontend (Flask) :8098
After=network-online.target jarwo-asr.service jarwo-tts.service
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/jarwo-assistant
ExecStart=/root/jarwo-assistant/.venv/bin/python server.py
Restart=always
RestartSec=5
Environment=JARWO_PORT=8098
Environment=JARWO_REMOTE_ASR_URL=http://127.0.0.1:8096/transcribe
Environment=JARWO_REMOTE_TTS_URL=http://127.0.0.1:8000/tts
StandardOutput=append:/root/jarwo-assistant/gateway.log
StandardError=append:/root/jarwo-assistant/gateway.log

[Install]
WantedBy=multi-user.target
